# templates/default/wlnodemanager.erb
#!/bin/sh
#
# wlnodemanager Oracle Weblogic NodeManager
#
# chkconfig:   345 85 15
# description: Oracle Weblogic service

### BEGIN INIT INFO
# Provides: wlnodemanager
# Required-Start: $network $local_fs
# Required-Stop:
# Should-Start:
# Should-Stop:
# Default-Start: 3 4 5
# Default-Stop: 0 1 2 6
# Short-Description: Oracle Weblogic service.
# Description: Starts and stops Oracle Weblogic NodeManager
### END INIT INFO

. /etc/rc.d/init.d/functions

DOMAIN_HOME="<%= node[:weblogic][:domain_home] %>/<%= node[:weblogic][:domain_name] %>"
NM_HOME="<%= node[:weblogic][:nodemanager_home] %>"
WLS_HOME="<%= node[:weblogic][:weblogic_home] %>"
DAEMON_USER="oracle"
PROCESS_STRING="^.*${DOMAIN_HOME}.*weblogic.NodeManager.*"

SERVICE_NAME=`/bin/basename $0`
LOCKFILE="${NM_HOME}/nodemanager.log.lck"

RETVAL=0

startNodeManager(){
  echo 'Starting WebLogic NodeManager.'
  /bin/su ${DAEMON_USER} -c "${WLS_HOME}/common/bin/wlst.sh ${DOMAIN_HOME}/startNodeManager.py &> ${NM_HOME}/startup.log"
}

stopNodeManager(){
  echo 'Stopping Weblogic NodeManager.'
  /bin/su ${DAEMON_USER} -c "${WLS_HOME}/common/bin/wlst.sh ${DOMAIN_HOME}/stopNodeManager.py &> ${NM_HOME}/shutdown.log"
}

start() {
  OLDPID=`/usr/bin/pgrep -f ${PROCESS_STRING}`
  if [ ! -z "$OLDPID" ]; then
    echo '${SERVICE_NAME} is already running (pid ${OLDPID}) !'
    exit
  fi
  startNodeManager
  RETVAL=$?
  echo
  return ${RETVAL}
}

stop() {
  stopNodeManager
  RETVAL=$?

  OLDPID=`/usr/bin/pgrep -f ${PROCESS_STRING}`
  if [ "$OLDPID" != "" ]; then
    kill -TERM ${OLDPID}
  else
    echo "${SERVICE_NAME} is stopped"
  fi
  return ${RETVAL}
}

restart() {
        stop
        sleep 10
        start
}

case "$1" in
  start)
        start
        ;;
  stop)
        stop
        ;;
  restart|force-reload|reload)
        restart
        ;;
  condrestart|try-restart)
        [ -f ${LOCKFILE} ] && restart
        ;;
  status)
        OLDPID=`/usr/bin/pgrep -f ${PROCESS_STRING}`
        if [ "${OLDPID}" != "" ]; then
            /bin/echo "${SERVICE_NAME} is running (pid: ${OLDPID})"
            RETVAL=0
        else
            /bin/echo "${SERVICE_NAME} is stopped"
            RETVAL=1
        fi
        ;;
  *)
        echo $"Usage: $0 {start|stop|status|restart|reload|force-reload|condrestart}"
        exit 1
esac

exit ${RETVAL}
